# services/rosnfr.service
[Unit]
Description=This is the ROS-NFR service.
[Service]
Type=simple
ExecStart=${HOME}/scripts/run_rosnfr.sh
User=root
[Install]