# services/rosnfr.service
[Unit]
Description=This is the ROS-NFR service.
[Service]
Type=simple
ExecStart=/opt/rosnfr/launch.sh
User=root
[Install]